# package/luci-app-airoha-flowsense/root/etc/init.d/npu-jitter
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

DAEMON=/usr/libexec/npu-jitter-daemon
PIDFILE=/var/run/npu-jitter.pid

start_service() {
	local target
	# Shipped configs before 2026-07 used an anonymous section ("config jitter"
	# with no name), only addressable as @jitter[0]; newer ones name it. Try
	# both so user-set targets work on either (thanks wolfgangrupprecht).
	target=$(uci -q get npu-monitor.jitter.target || \
		 uci -q get npu-monitor.@jitter[0].target)
	# No configured target: measure against the default gateway (router mode:
	# the ISP first hop; AP mode: the main router) - always reachable, no
	# hardcoded external host. Last resort keeps the old external default.
	[ -n "$target" ] || \
		target=$(ip -4 route show default 2>/dev/null | awk '/via/ {print $3; exit}')
	[ -n "$target" ] || target="1.1.1.1"

	procd_open_instance
	procd_set_param command "$DAEMON" "$target"
	procd_set_param pidfile "$PIDFILE"
	procd_set_param respawn 30 5 0
	procd_set_param stdout 0
	procd_set_param stderr 0
	procd_close_instance
}

stop_service() {
	rm -f /tmp/npu-jitter.json /tmp/npu-jitter.json.tmp
}

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "npu-monitor"
}
